# app/install/power_snitch_web.service
[Unit]
Description=Power Snitch Web Interface
After=network.target power_snitch.service
Wants=power_snitch.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/power_snitch
Environment=PYTHONUNBUFFERED=1
Environment=FLASK_ENV=production
Environment=FLASK_APP=web/app.py
Environment=PYTHONPATH=/opt/power_snitch
ExecStart=/usr/bin/python3 /opt/power_snitch/web/app.py
Restart=always
RestartSec=10
StandardOutput=append:/var/log/power_snitch/web_app.service.log
StandardError=append:/var/log/power_snitch/web_app.service.error.log
ProtectSystem=full
ProtectHome=true
NoNewPrivileges=true
PrivateTmp=true
SystemCallArchitectures=native
CapabilityBoundingSet=CAP_NET_BIND_SERVICE

[Install]
WantedBy=multi-user.target 