# power_snitch.service
[Unit]
Description=Power Snitch UPS Monitoring Service
After=network.target nut-server.service
Wants=nut-server.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/power_snitch
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python3 /opt/power_snitch/power_snitch.py
Restart=always
RestartSec=10
StandardOutput=append:/var/log/power_snitch/power_snitch.service.log
StandardError=append:/var/log/power_snitch/power_snitch.service.error.log

[Install]
WantedBy=multi-user.target 