# archive/app/install/power_snitch.service
[Unit]
Description=Power Snitch UPS Monitoring Service
After=network.target nut-server.service
Wants=nut-server.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/power_snitch
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=/opt/power_snitch
Environment=FLASK_ENV=production
Environment=FLASK_APP=power_snitch.py
ExecStart=/usr/bin/python3 /opt/power_snitch/power_snitch.py
Restart=always
RestartSec=10
StandardOutput=append:/var/log/power_snitch/power_snitch.service.log
StandardError=append:/var/log/power_snitch/power_snitch.service.error.log
ProtectSystem=full
ProtectHome=true
NoNewPrivileges=true
PrivateTmp=true
SystemCallArchitectures=native
CapabilityBoundingSet=CAP_NET_BIND_SERVICE

[Install]
WantedBy=multi-user.target 